# private_dot_config/systemd/user/registry-docker.io.socket
[Unit]
Description=Registry Docker.io socket
Wants=network-online.target
After=network-online.target

[Socket]
ListenStream=5000

[Install]
WantedBy=sockets.target
